# packaging/systemd/smtpserver.service
[Unit]
Description=Ship Shape Sharpening SMTP Server v0.4.0
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=smtpserver
Group=smtpserver
EnvironmentFile=/etc/smtpserver/smtpserver.env
WorkingDirectory=/opt/smtpserver
ExecStart=/usr/local/bin/smtpserver
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/smtpserver /var/spool/smtpserver /etc/smtpserver
RuntimeDirectory=smtpserver
StateDirectory=smtpserver
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
